# vault-agent.service
[Unit]
Description=HashiCorp Vault Agent Service

[Service]
ExecStart=./restart-agent.sh

[Install]
WantedBy=multi-user.target